# beautybot railway/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=50MB
pidfile=/tmp/supervisord.pid

[program:backend]
command=node backend/src/index.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3001",DB_PATH="/app/backend/data/beauty_salon.db",UPLOADS_PATH="/app/backend/uploads"

[program:bot]
command=node bot/src/index.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3000",BACKEND_URL="http://localhost:3001"
